# systemd/openkoutsi-wahoo-bridge@.service
[Unit]
Description=openkoutsi Wahoo bridge
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/projects/openkoutsi-backend/wahoo_bridge
EnvironmentFile=/home/%i/projects/openkoutsi-backend/wahoo_bridge/.env
ExecStart=/home/%i/.local/bin/uv run uvicorn main:app --host 0.0.0.0 --port 8085
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
